# app/hermesd.service
[Unit]
Description = Hermes - Servicio destinado a la actualización automática de servicios Linux
After = multi-user.target

[Service]
Type = simple
Restart = on-failure
ExecStart = /usr/bin/python3 /hermesd/hermes.py

[Install]
WantedBy=multi-user.target