# app/hermesd.service
[Unit]
Description = Hermes - Linux service designed to automate the management of installed services
Requires = hermes-dashboard.service
After = hermes-dashboard.service

[Service]
ExecStart = python3 /hermesd/hermes.py &

[Install]
WantedBy = multi-user.target
