# app/hermesd.service
[Unit]
Description = Hermes - Linux service designed to automate the management of installed services
After = multi-user.target


[Service]
ExecStart = python3 /hermesd/hermes.py &
ExecStop = python3 /hermesd/hermes.py stop

[Install]
WantedBy = multi-user.target
