# app/hermesd.service
[Unit]
Description = Hermes - Linux service designed to automate the management of installed services
After = multi-user.target

[Service]
Type = simple
Restart = on-failure
ExecStart = /usr/bin/python3 /hermesd/hermes.py

[Install]
WantedBy=multi-user.target