# app/hermesd.service
[Unit]
Description = Hermes - Linux service designed to automate the management of installed services
After = multi-user.target

[Service]
Type=forking
ExecStart= /hermesd/start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target