# app/hermesd.service
[Unit]
Description = Hermes - Linux service designed to automate the management of installed services
After = multi-user.target
Requires = hermesd-dashboard.service

[Service]
ExecStart = python3 /hermesd/hermes.py

[Install]
WantedBy = multi-user.target